# scripts/systemd/stocks-backend.service
[Unit]
Description=Stocks Backend (FastAPI on :8000)
After=network.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}/backend
ExecStart={{PROJECT_ROOT}}/backend/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
StandardOutput=append:{{PROJECT_ROOT}}/backend/server.log
StandardError=append:{{PROJECT_ROOT}}/backend/server.log
KillMode=mixed

[Install]
WantedBy=default.target
